# scripts/systemd/vm-monitor-report.service
[Unit]
Description=VM Monitoring Health Report (Periodic Status)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vm-monitor.sh --mode report
StandardOutput=journal
StandardError=journal
